# app/conf.d/clone_repo.conf

[program:git_clone]
command=/app/bin/clone_repo.sh
autostart=true
autorestart=false
startsecs=0
user=vscode
