# app/conf.d/x11vnc.conf
[program:x11vnc]
#command=x11vnc -forever -shared
command=x11vnc -display :0 -forever -shared -nopw -listen 0.0.0.0 -rfbport 5900 -xkb
autorestart=true
user=root
